# service/file-uploader.service
[Unit]
Description=File upload server
After=syslog.target network.target

[Service]
ExecStart=/xxx/xxx/uploader
User=www-data

[Install]
WantedBy=multi-user.target
